# supervisord.conf
; supervisord.conf — manages uvicorn and the Celery worker inside one container.
;
; Both programs log to the container's stdout/stderr so HF Spaces (and any
; container orchestrator) can collect logs with `docker logs`.

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ── uvicorn (FastAPI) ─────────────────────────────────────────────────────────
[program:uvicorn]
command=uvicorn app.main:app --host 0.0.0.0 --port 7860 --workers 1
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"

; ── Celery worker ─────────────────────────────────────────────────────────────
[program:celery]
command=celery -A app.worker.celery_app worker --loglevel=info --concurrency=2
directory=/app
autostart=true
autorestart=true
startsecs=10
startretries=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"
